# wordpress/docker/include/supervisord.ini
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "exec httpd -D FOREGROUND -f /etc/apache2/httpd.conf"
autostart=yes
autorestart=true

[program:php-fpm]
command=/bin/bash -c "exec php-fpm"
autostart=yes
autorestart=true